# infrastructure/alpha-heartbeat.timer
[Unit]
Description=Alpha Heartbeat Timer - Periodic activation for Project Beta

[Timer]
# Midnight to 4 AM, every 15 minutes (17 heartbeats)
# First night: December 8-9, 2025
# Hours 00-03: every 15 minutes
OnCalendar=*-*-* 00..03:00,15,30,45:00
# Hour 04: just the top of the hour (final heartbeat)
OnCalendar=*-*-* 04:00:00

# Don't run missed heartbeats if the Pi was asleep
Persistent=false

# Add some randomness to avoid exact-second firing
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
